# mongo/mongo.sv.conf
[program:mongo]
command=/opt/mongodb/bin/mongod -f /etc/mongod.conf
redirect_stderr=true
